# protected/extensions/crontab/crontabs/cron_requests
40 17 25 3 * php /var/www/razbor/yiic.php request broken --id=1
