# protected/extensions/crontab/crontabs/cron_req_20
28 14 25 3 * php /var/www/razbor/protected/yiic.php request broken --id=20
